# rootfs/etc/s6-overlay/s6-rc.d/lsyncd/finish
#!/command/with-contenv sh

if [ "$LSYNCD" = "true" ]; then
    echo "+正在停止 Lsyncd 守护进程..."
else
    exit 0
fi