# rootfs/etc/s6-overlay/s6-rc.d/init/up
#!/command/execlineb

foreground {
    /etc/s6-overlay/s6-rc.d/init/init.sh
}
/etc/s6-overlay/s6-rc.d/init/banner.sh