# rootfs/etc/s6-overlay/s6-rc.d/sshd/finish
#!/command/with-contenv sh

if [ "$SSH" = "true" ]; then
    echo "+正在停止 SSHD 守护进程..."
    rm -f /conf/.ssh/known_hosts*
else
    exit 0
fi