# rootfs/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/with-contenv sh

if [ "$SSH" = "true" ]; then
    echo "+正在启动 SSHD 守护进程..."
    exec /usr/sbin/sshd -D -e
else
    s6-svc -Od .
    exit 0
fi